# node01/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:zookeeper]
command=/opt/message/zookeeper_start.sh

[program:zookeeper]
command=/opt/init/init.sh
autorestart=false
